# debian/parch-util.cron.d
15 5,17 * * *  wdb /usr/bin/parch_locationforecastLoad -uwdb -hlocalhost
